# deployment/gunicorn.service
[Unit]
Description=gunicorn daemon for Elogforlinux (elogbookagu)
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/projects/Elogforlinux/elogbookagu
ExecStart=/home/ubuntu/projects/myenv/bin/gunicorn \
  --workers 3 \
  --bind 127.0.0.1:8001 \
  --timeout 30 \
  --log-level info \
  --access-logfile /home/ubuntu/projects/Elogforlinux/elogbookagu/gunicorn-access.log \
  --error-logfile /home/ubuntu/projects/Elogforlinux/elogbookagu/gunicorn-error.log \
  elogbookagu.wsgi:application

Environment=PATH=/home/ubuntu/projects/myenv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
Environment=PYTHONUNBUFFERED=1
Environment=DJANGO_SETTINGS_MODULE=elogbookagu.settings

Restart=always
RestartSec=5
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
